# ansible/autoload/rc.local
#!/bin/sh -e
/usr/local/bin/updater_start.sh
exit 0
